# htdocs/crontab/root
1 0 * * * /usr/bin/bash /data/www/website/bak.sh website sql
